# deploy/star-office-bridge.service
[Unit]
Description=Star Office Bridge - Bun HTTP server for Zellij/Claude event routing
After=network.target

[Service]
Type=simple
ExecStart=/root/.bun/bin/bun run src/index.ts
WorkingDirectory=/root/free-code-star-office-zellij-lab
Environment=BRIDGE_HOST=127.0.0.1
Environment=BRIDGE_PORT=4317
Environment=BRIDGE_DRY_RUN=true
Environment=ZELLIJ_WEB_URL=https://term.karldigi.dev/main
Environment=ZELLIJ_SESSION_NAME=main
Environment=PATH=/root/.bun/bin:/usr/local/bin:/usr/bin:/bin
# Secrets loaded from local env file (not tracked in git)
EnvironmentFile=-/etc/star-office-bridge.env

# Restart policy
Restart=on-failure
RestartSec=3
UMask=0077

# Graceful shutdown: give 30s for server.stop() + SSE drain before SIGKILL
TimeoutStopSec=30

# High connection count for SSE clients
LimitNOFILE=65535

# === MEMORY PROTECTION ===
# Soft limit: throttle when exceeding 256MB (kernel applies pressure)
MemoryHigh=256M
# Hard limit: kill at 512MB to prevent runaway memory from crashing the host
MemoryMax=512M
# Protected baseline: 64MB guaranteed before reclaim
MemoryLow=64M
# Cap swap usage to prevent silent memory creep
MemorySwapMax=128M
# On OOM kill: restart the service (don't leave it dead)
OOMPolicy=stop

# Logging: journal integration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bridge

# === SECURITY HARDENING ===

# Bridge now proxies zellij action commands which need /tmp/zellij-0/ socket access
# PrivateTmp=yes

# ProtectHome: DISABLED -- zellij action commands need write access to /root/.cache/zellij/
# ProtectHome=read-only

# ProtectSystem: full (read-only /usr, /boot, /etc, /var) — allows /tmp and /root writes
# strict mode breaks zellij action commands which write to /root/.cache/zellij/
ProtectSystem=full

NoNewPrivileges=yes
PrivateDevices=yes
ProtectProc=invisible
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictNamespaces=yes
# Bun is a JS runtime -- MDWX may break JIT. Test carefully.
# If Bun crashes with SIGILL, comment out MemoryDenyWriteExecute
# MemoryDenyWriteExecute=yes
CapabilityBoundingSet=~CAP_SYS_PTRACE ~CAP_SYS_ADMIN ~CAP_SYS_RAWIO ~CAP_SYS_MODULE

[Install]
WantedBy=multi-user.target
