# deploy/zellij-subscribe-stream.service
[Unit]
Description=Zellij Subscribe Stream - streams pane content via subscribe --format json
After=star-office-bridge.service
Wants=star-office-bridge.service

[Service]
Type=simple
ExecStartPre=/bin/sleep 3
ExecStart=/root/free-code-star-office-zellij-lab/scripts/zellij-subscribe-stream.sh
WorkingDirectory=/root/free-code-star-office-zellij-lab

# Environment — no secrets here (loaded via LoadCredential)
Environment=BRIDGE_URL=http://127.0.0.1:4317
Environment=ZELLIJ_SESSION_NAME=main
Environment=CONTENT_DEBOUNCE=3
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=HOME=/root
Environment=XDG_RUNTIME_DIR=/run/user/0

# Secret: loaded from credstore, not visible via systemctl show or /proc/environ
LoadCredential=bridge-secret:/etc/credstore/star-office-bridge.secret

# Restart policy
Restart=on-failure
RestartSec=5
UMask=0077

# Logging: journal integration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zellij-subscribe

# === SECURITY HARDENING ===

# PrivateTmp: DISABLED -- Zellij stores session sockets at /tmp/zellij-0/
# PrivateTmp=yes

# ProtectHome: DISABLED -- Zellij needs write access to /root/.cache/zellij/
# ProtectHome=read-only

# ProtectSystem: full (read-only /usr, /boot, /etc, /var) — allows /tmp and /root writes
ProtectSystem=full

# Prevent privilege escalation
NoNewPrivileges=yes

# Limit device access
PrivateDevices=yes

# Restrict networking to localhost + Unix sockets
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Process visibility
ProtectProc=invisible

# Kernel protection
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Capability and syscall restrictions
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictNamespaces=yes
MemoryDenyWriteExecute=yes
CapabilityBoundingSet=~CAP_SYS_PTRACE ~CAP_SYS_ADMIN ~CAP_SYS_RAWIO ~CAP_SYS_MODULE

[Install]
WantedBy=multi-user.target
